# systemd/system/snap-google\x2dcloud\x2dsdk-124.mount
[Unit]
Description=Mount unit for google-cloud-sdk, revision 124
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/google-cloud-sdk_124.snap
Where=/snap/google-cloud-sdk/124
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
